# piper-service/piper-tts.service
# Systemd service file for Piper TTS Docker service
# Install to: /etc/systemd/system/piper-tts.service
# 
# Installation:
#   sudo cp piper-tts.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable piper-tts.service
#   sudo systemctl start piper-tts.service

[Unit]
Description=Piper TTS Service (Docker)
Documentation=https://github.com/rhasspy/piper
After=docker.service
Requires=docker.service
PartOf=blindr-services.target

[Service]
Type=notify
Restart=always
RestartSec=10
TimeoutStartSec=300
TimeoutStopSec=30

# User and group for Docker operations
User=root
Group=docker

# Working directory
WorkingDirectory=/home/travis/services/piper-service

# Environment variables
Environment=DOCKER_BUILDKIT=1
Environment=COMPOSE_DOCKER_CLI_BUILD=1

# Pre-start checks
ExecStartPre=-/usr/bin/docker stop piper-tts-service
ExecStartPre=-/usr/bin/docker rm piper-tts-service
ExecStartPre=/usr/bin/docker pull piper-tts-service:latest

# Start the service
ExecStart=/usr/bin/docker run \
    --name piper-tts-service \
    --rm \
    -p 9001:9001 \
    --memory=2g \
    --cpus=1.0 \
    --health-interval=30s \
    --health-timeout=15s \
    --health-start-period=45s \
    --health-retries=3 \
    -v piper-models:/app/models/piper \
    -v piper-logs:/app/logs \
    -v piper-tmp:/app/tmp \
    -e PIPER_PORT=9001 \
    -e PIPER_HOST=0.0.0.0 \
    -e LOG_LEVEL=INFO \
    -e PIPER_DEFAULT_VOICE=en_US-lessac-medium \
    -e PIPER_MAX_TEXT_LENGTH=10000 \
    -e PIPER_MAX_CONCURRENT=10 \
    -e PIPER_MODEL_CACHE_SIZE=5 \
    piper-tts-service:latest

# Stop the service
ExecStop=/usr/bin/docker stop -t 30 piper-tts-service

# Post-stop cleanup
ExecStopPost=-/usr/bin/docker rm piper-tts-service

# Reload configuration
ExecReload=/usr/bin/docker restart piper-tts-service

# Health monitoring
WatchdogSec=60

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=piper-tts-service

[Install]
WantedBy=multi-user.target
WantedBy=blindr-services.target

# Service management commands:
#
# Status:  sudo systemctl status piper-tts.service
# Logs:    sudo journalctl -u piper-tts.service -f
# Start:   sudo systemctl start piper-tts.service  
# Stop:    sudo systemctl stop piper-tts.service
# Restart: sudo systemctl restart piper-tts.service
# Enable:  sudo systemctl enable piper-tts.service
# Disable: sudo systemctl disable piper-tts.service